# templates/default/graylog_web.conf.erb
# Ampel Service

description     "Graylog2 Webservice"
author          "Daniel Hahn"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 2 5

env HOME=<%= node.graylog2.basedir -%>/web

script
  cd $HOME
  exec sudo -u <%= node.graylog2.web_interface_user -%> sh -c "<%= node.graylog2.basedir -%>/web/bin/graylog2-web-interface -Dhttp.port=<%= node.graylog2.web_interface_port -%>"
end script
